# deploy/systemd/agent-telegram.service
[Unit]
Description=Leo Agent Telegram Poller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/leoprojects/agent-cli-v4-fresh
EnvironmentFile=/opt/leoprojects/agent-cli-v4-fresh/.env
ExecStart=/opt/leoprojects/agent-cli-v4-fresh/.venv/bin/python telegram_poll.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
